# chroot/lib/udev/rules.d/40-wipe.rules
ACTION=="add", SUBSYSTEM=="block", RUN+="/usr/local/wipe.sh"
